# extras/supervisor/igd.conf
[program:igd]
user=username
directory=/home/username/MELGA/IGD
command=/home/username/MELGA/IGD/igd
autorestart=true
redirect_stderr=true
