# extras/supervisor/igd.conf
[program:igd]
directory=/home/dgomes/MELGA/IGD
command=/home/dgomes/MELGA/IGD/igd
autorestart=true
