# extras/supervisor/mqttwarn.conf
[program:mqttwarn]
directory = /opt/mqttwarn
command = /opt/mqttwarn/mqttwarn.py
user = root 
environment= MQTTWARNINI="/opt/mqttwarn/mqttwarn.ini"
