# extras/supervisor/xively.conf
[program:xively]
user=username
directory=/home/username/MELGA/xively
command=/home/username/MELGA/xively/xively 
autorestart=true
redirect_stderr=true
